# supervisor/activemq.conf
[program:activemq]
command=/opt/apache-activemq-5.9.0/bin/activemq start
autorestart=true